# com.user.timemachine-auto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- LaunchAgent commands (run in Terminal):
         unload: launchctl bootout "gui/$(id -u)"/com.user.timemachine-auto 2>/dev/null || true
         load:   launchctl bootstrap "gui/$(id -u)" "$HOME/Library/LaunchAgents/com.user.timemachine-auto.plist" && launchctl enable "gui/$(id -u)"/com.user.timemachine-auto
         optional verification that it is running: launchctl enable "gui/$(id -u)"/com.user.timemachine-auto
    -->
    <!-- Unique launchd job identifier. -->
    <key>Label</key>
    <string>com.user.timemachine-auto</string>

    <!-- Command and arguments launchd executes. -->
    <key>ProgramArguments</key>
    <array>
        <!-- Shell executable. -->
        <string>/bin/bash</string>
        <!-- Run command string mode. -->
        <string>-c</string>
        <!-- Script path in user Library; update if you deploy elsewhere. -->
        <string>exec "$HOME/Library/Scripts/timemachine-auto.sh"</string>
    </array>

    <!-- Trigger whenever any filesystem is mounted. -->
    <key>StartOnMount</key>
    <true/>

    <!-- Run once when the agent is loaded at login/session start. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restrict to GUI user session so notifications work. -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>
</dict>
</plist>
