# com.user.timemachine-auto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.timemachine-auto</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME_HERE/Library/Scripts/timemachine-auto.sh</string>
    </array>
    <key>StartOnMount</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME_HERE/Library/Logs/AutoTMLogs/tm-launchd-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME_HERE/Library/Logs/AutoTMLogs/tm-launchd-err.log</string>
</dict>
</plist>